# infra/scripts/systemd/ffmpeg-stream.service
[Unit]
Description=FFmpeg Stream Capture + YouTube RTMP
After=minecraft-client.service avatar-writer.service
Requires=avatar-writer.service

[Service]
Type=simple
User=ubuntu
Environment=DISPLAY=:99
EnvironmentFile=/home/ubuntu/ai-minecraft/.env
ExecStart=/usr/bin/ffmpeg \
  -f x11grab -video_size 1280x720 -framerate 24 -i :99 \
  -f pulse -i combined_sink.monitor \
  -f rawvideo -pixel_format rgba -video_size 300x400 -framerate 5 -i /tmp/ai-minecraft-avatar.pipe \
  -filter_complex "[0:v][2:v]overlay=W-w-20:H-h-20:format=auto[out]" \
  -map "[out]" -map 1:a \
  -c:v libx264 -preset ultrafast -tune zerolatency -b:v 3000k -maxrate 3500k -bufsize 6000k -r 24 \
  -c:a aac -b:a 128k -ar 44100 \
  -f flv rtmp://a.rtmp.youtube.com/live2/${YOUTUBE_STREAM_KEY}
Restart=on-failure
RestartSec=5
CPUAffinity=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
