# infra/scripts/systemd/orchestrator.service
[Unit]
Description=AI Minecraft Orchestrator
After=minecraft-server.service voicevox.service
Wants=minecraft-server.service voicevox.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ai-minecraft
EnvironmentFile=/home/ubuntu/ai-minecraft/.env
Environment=NODE_ENV=production
ExecStartPre=/bin/bash /home/ubuntu/.config/pulse/set-volumes.sh
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=15
CPUAffinity=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
